# deploy/supervisord.conf
#********************************************************************
# 项目名称：DevOps
# 文件名称：supervisord.conf
# 创建时间：2026-04-14 22:33:28
#
# 系统用户：jerion
# 作　　者：Jerion
# 联系邮箱：416685476@qq.com
# 功能描述：supervisord 多进程管理配置，管理 Nginx 和 Go 后端服务
#********************************************************************

[supervisord]
nodaemon=true
user=root
logfile=/app/logs/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
stdout_logfile=/app/logs/nginx-stdout.log
stderr_logfile=/app/logs/nginx-stderr.log
priority=10

[program:devops-server]
command=/app/devops
autostart=true
autorestart=true
stdout_logfile=/app/logs/devops-stdout.log
stderr_logfile=/app/logs/devops-stderr.log
priority=20
